# systemd/user/backup_packages.service
[Unit]
Description=Backup packages list to file in .config

[Service]
Type=oneshot
ExecStart=/home/Samuel/.config/backup-packages.sh
